# packages/3rdparty/service/aml-vnc/source/system.d/service.system.aml-vnc.service
[Unit]
Description = Amlogic VNC Server
After = network-online.target graphical.target

[Service]
Type = simple
ExecStartPre = /bin/sleep 5
ExecStart = /bin/sh -c "exec sh /storage/.kodi/addons/service.system.aml-vnc/bin/aml-vnc.start"
TimeoutStopSec = 20
Restart = always
RestartSec = 0

[Install]
WantedBy = kodi.target
Alias = aml-vnc.service
