# packages/3rdparty/service/aml-vnc/source/system.d/service.system.aml-vnc.service
[Unit]
Description=Amlogic VNC server
After=kodi.service network-online.target
Requires=kodi.service network-online.target

[Service]
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.system.aml-vnc/bin/aml-vnc.start"
TimeoutStopSec=1
Restart=always
RestartSec=2
StartLimitInterval=0

[Install]
WantedBy=kodi.target
