# packages/3rdparty/libraries/makemkv/source/profile.d/makemkv.profile
MAKEMKV_PATH=/storage/.kodi/addons/lib.multimedia.makemkv
LIBBDPLUS_PATH="$MAKEMKV_PATH/lib/libmmbd"
LIBAACS_PATH="$MAKEMKV_PATH/lib/libmmbd"
MAKEMKVCON="$MAKEMKV_PATH/bin/makemkvcon.bin"
export LIBBDPLUS_PATH
export LIBAACS_PATH
export MAKEMKVCON
