# packages/3rdparty/service/makemkv/source/profile.d/makemkv.profile
ADDON_DIR=/storage/.kodi/addons/service.makemkv
LIBBDPLUS_PATH="$ADDON_DIR/lib/libmmbd"
LIBAACS_PATH="$ADDON_DIR/lib/libmmbd"
MAKEMKVCON="$ADDON_DIR/bin/makemkvcon"
export LIBBDPLUS_PATH
export LIBAACS_PATH
export MAKEMKVCON
